# systemd/user/local-inference.service.d/rpc-node.example.conf
[Service]
# Example RPC-node override. Copy to:
#   ~/.config/systemd/user/local-inference.service.d/rpc-node.conf
# and adjust paths locally; do not commit the live override.
Environment=LOCAL_INFERENCE_DATA_ROOT=/mnt/data/local-inference
Environment=LOCAL_INFERENCE_VENV=/mnt/data/local-inference/.venv-local-inference
Environment=QWEN_TEXT_EMBEDDING_MODEL_DIR=/mnt/data/local-inference/models/Qwen3-Embedding-8B
Environment=QWEN_RERANKER_MODEL_DIR=/mnt/data/local-inference/models/Qwen3-Reranker-8B
Environment=QWEN_VL_EMBEDDING_MODEL_DIR=/mnt/data/local-inference/models/Qwen3-VL-Embedding-8B
Environment=QWEN_VL_GENERATIVE_MODEL_DIR=/mnt/data/local-inference/models/Qwen3-VL-8B-Instruct
